# chronystats.cron.default
*/5 * * * * root /root/pythons/chronylog/chronylog.py > /dev/null 2>&1 
